# packages/lakka/lakka_tools/retroflag_picase_safeshutdown/system.d/retroflag_picase_safeshutdown.service
[Unit]
Description=RETROFLAG Pi CASE series - safe shutdown service
ConditionPathExists=|/storage/.cache/services/safeshutdown.conf

[Service]
ExecStart=/usr/bin/retroflag_picase_safeshutdown.py
TimeoutStopSec=5
Restart=always
RestartSec=5
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
